# dist/systemd/etc/systemd/system/tac_plus.service
[Unit]
Description=Shrubbery TACACS+ Server
After=network.target auditd.service

[Service]
ExecStart=/usr/sbin/tac_plus -C /etc/tacacs/tac_plus.conf
KillMode=process
Restart=on-failure
